# udevil/DEBIAN/prerm
#! /bin/sh
set -e

if which systemctl > /dev/null;then
	set +e
	systemctl stop devmon.service > /dev/null 2>&1
	systemctl disable devmon.service
	set -e
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
